# deploy/style-ai-site-worker.service
[Unit]
Description=style-ai-site queue worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/style-ai-site/server
EnvironmentFile=/opt/style-ai-site/server/.env
ExecStart=/usr/bin/npm run worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
